# owlgorithm.cron
30 9 * * * /Users/jonamar/Documents/owlgorithm/run_tracker.sh 