# owlgorithm.cron
# Duolingo Progress Tracker - Runs 4 times daily
# 6:00 AM - Morning check
0 6 * * * cd /Users/jonamar/Documents/owlgorithm && /Users/jonamar/Documents/owlgorithm/duolingo_env/bin/python /Users/jonamar/Documents/owlgorithm/scripts/daily_update.py >> /Users/jonamar/Documents/owlgorithm/logs/tracker.log 2>&1

# 12:00 PM - Midday check  
0 12 * * * cd /Users/jonamar/Documents/owlgorithm && /Users/jonamar/Documents/owlgorithm/duolingo_env/bin/python /Users/jonamar/Documents/owlgorithm/scripts/daily_update.py >> /Users/jonamar/Documents/owlgorithm/logs/tracker.log 2>&1

# 5:00 PM - Evening check
0 17 * * * cd /Users/jonamar/Documents/owlgorithm && /Users/jonamar/Documents/owlgorithm/duolingo_env/bin/python /Users/jonamar/Documents/owlgorithm/scripts/daily_update.py >> /Users/jonamar/Documents/owlgorithm/logs/tracker.log 2>&1

# 10:00 PM - Night check
0 22 * * * cd /Users/jonamar/Documents/owlgorithm && /Users/jonamar/Documents/owlgorithm/duolingo_env/bin/python /Users/jonamar/Documents/owlgorithm/scripts/daily_update.py >> /Users/jonamar/Documents/owlgorithm/logs/tracker.log 2>&1 